# deploy/myclaw-api.service
[Unit]
Description=MyClaw API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/myclaw/apps/api
ExecStart=/usr/bin/pnpm start
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=2222
EnvironmentFile=/home/ubuntu/myclaw/apps/api/.env

[Install]
WantedBy=multi-user.target
